# apps/ethos/templates/systemd/ethos-runall.service.tmpl
[Unit]
Description=Ethos All Services (gateway + web API)
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User={{ETHOS_USER}}
Environment=HOME={{ETHOS_HOME}}
Environment=ETHOS_MANAGED=1
EnvironmentFile=-{{ETHOS_HOME}}/.ethos/.env
ExecStart={{ETHOS_BINARY}} run-all
Restart=on-failure
RestartSec=5
WatchdogSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
